# sleep.conf
description "sleep command"
author  "shigure.onishi@hacomono.co.jp"

start on runlevel [2345]
stop on runlevel [016]

chdir /
exec sleep 100000000
respawn